# rocketcher-web.service
[Unit]
Description=rocketcher web service

[Service]
Type=simple
ExecStart=<al-zal-ddak>
Restart=on-failure

[Install]
WantedBy=multi-user.target